# build/vcpkg_installed/vcpkg/info/dbus_1.15.8_x64-linux.list
x64-linux/
x64-linux/debug/
x64-linux/debug/lib/
x64-linux/debug/lib/dbus-1.0/
x64-linux/debug/lib/dbus-1.0/include/
x64-linux/debug/lib/dbus-1.0/include/dbus/
x64-linux/debug/lib/dbus-1.0/include/dbus/dbus-arch-deps.h
x64-linux/debug/lib/libdbus-1.so
x64-linux/debug/lib/libdbus-1.so.3
x64-linux/debug/lib/libdbus-1.so.3.38.0
x64-linux/debug/lib/pkgconfig/
x64-linux/debug/lib/pkgconfig/dbus-1.pc
x64-linux/debug/lib/systemd/
x64-linux/debug/lib/systemd/system/
x64-linux/debug/lib/systemd/system/dbus.service
x64-linux/debug/lib/systemd/system/dbus.socket
x64-linux/debug/lib/systemd/system/multi-user.target.wants/
x64-linux/debug/lib/systemd/system/multi-user.target.wants/dbus.service
x64-linux/debug/lib/systemd/system/sockets.target.wants/
x64-linux/debug/lib/systemd/system/sockets.target.wants/dbus.socket
x64-linux/debug/libexec/
x64-linux/debug/libexec/dbus-daemon-launch-helper
x64-linux/include/
x64-linux/include/dbus-1.0/
x64-linux/include/dbus-1.0/dbus/
x64-linux/include/dbus-1.0/dbus/dbus-address.h
x64-linux/include/dbus-1.0/dbus/dbus-bus.h
x64-linux/include/dbus-1.0/dbus/dbus-connection.h
x64-linux/include/dbus-1.0/dbus/dbus-errors.h
x64-linux/include/dbus-1.0/dbus/dbus-macros.h
x64-linux/include/dbus-1.0/dbus/dbus-memory.h
x64-linux/include/dbus-1.0/dbus/dbus-message.h
x64-linux/include/dbus-1.0/dbus/dbus-misc.h
x64-linux/include/dbus-1.0/dbus/dbus-pending-call.h
x64-linux/include/dbus-1.0/dbus/dbus-protocol.h
x64-linux/include/dbus-1.0/dbus/dbus-server.h
x64-linux/include/dbus-1.0/dbus/dbus-shared.h
x64-linux/include/dbus-1.0/dbus/dbus-signature.h
x64-linux/include/dbus-1.0/dbus/dbus-syntax.h
x64-linux/include/dbus-1.0/dbus/dbus-threads.h
x64-linux/include/dbus-1.0/dbus/dbus-types.h
x64-linux/include/dbus-1.0/dbus/dbus.h
x64-linux/lib/
x64-linux/lib/dbus-1.0/
x64-linux/lib/dbus-1.0/include/
x64-linux/lib/dbus-1.0/include/dbus/
x64-linux/lib/dbus-1.0/include/dbus/dbus-arch-deps.h
x64-linux/lib/libdbus-1.so
x64-linux/lib/libdbus-1.so.3
x64-linux/lib/libdbus-1.so.3.38.0
x64-linux/lib/pkgconfig/
x64-linux/lib/pkgconfig/dbus-1.pc
x64-linux/lib/systemd/
x64-linux/lib/systemd/system/
x64-linux/lib/systemd/system/dbus.service
x64-linux/lib/systemd/system/dbus.socket
x64-linux/lib/systemd/system/multi-user.target.wants/
x64-linux/lib/systemd/system/multi-user.target.wants/dbus.service
x64-linux/lib/systemd/system/sockets.target.wants/
x64-linux/lib/systemd/system/sockets.target.wants/dbus.socket
x64-linux/libexec/
x64-linux/libexec/dbus-daemon-launch-helper
x64-linux/share/
x64-linux/share/DBus1/
x64-linux/share/DBus1/DBus1Config.cmake
x64-linux/share/DBus1/DBus1ConfigVersion.cmake
x64-linux/share/dbus-1/
x64-linux/share/dbus-1/session.conf
x64-linux/share/dbus/
x64-linux/share/dbus/copyright
x64-linux/share/dbus/vcpkg.spdx.json
x64-linux/share/dbus/vcpkg_abi_info.txt
x64-linux/tools/
x64-linux/tools/dbus/
x64-linux/tools/dbus/dbus-cleanup-sockets
x64-linux/tools/dbus/dbus-daemon
x64-linux/tools/dbus/dbus-launch
x64-linux/tools/dbus/dbus-monitor
x64-linux/tools/dbus/dbus-run-session
x64-linux/tools/dbus/dbus-send
x64-linux/tools/dbus/dbus-test-tool
x64-linux/tools/dbus/dbus-update-activation-environment
x64-linux/tools/dbus/dbus-uuidgen
